# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/etc/init.d/detbutton
#!/bin/sh /etc/rc.common
START=75

start() {
        # Start detbutton
        /sbin/button_detecte &
}

stop() {
        killall button_detecte
}

restart() {
        stop
        start
}
