# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/etc/init.d/wlanlog
#!/bin/sh /etc/rc.common
START=26

start() {
	# start wlanlog
	/usr/bin/wlanlog
}

stop() {
	killall /usr/bin/wlanlog
}

restart() {
	stop
	start
}
