# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/atheros-aquila-madwifi-11n-scripts.list
/
/usr/
/usr/lib/
/usr/lib/functions.iptv
/usr/lib/stainfo.awk
/lib/
/lib/wifi/
/lib/wifi/atheros11n.sh
/etc/
/etc/ath
/etc/dni-wifi-config
/etc/ath.orig/
/etc/ath.orig/wpscli_sta
/etc/ath.orig/killVAP
/etc/ath.orig/makeVAP
/etc/ath.orig/EAP.ap_bss
/etc/ath.orig/PSK.sta
/etc/ath.orig/RADIUS_MAC.ap_bss
/etc/ath.orig/board.conf
/etc/ath.orig/apdown
/etc/ath.orig/wifi_toggle
/etc/ath.orig/wpscli_ap
/etc/ath.orig/activateVAP
/etc/ath.orig/apcfg
/etc/ath.orig/WSC_sta.conf
/etc/ath.orig/PSK.ap_bss
/etc/ath.orig/WSC.conf
/etc/ath.orig/WSC_WEP.conf
/etc/ath.orig/wifi_schedule
/etc/ath.orig/hostapd0.7.0_conf/
/etc/ath.orig/hostapd0.7.0_conf/EAP.ap_bss
/etc/ath.orig/hostapd0.7.0_conf/PSK.ap_bss
/etc/ath.orig/hostapd0.7.0_conf/WSC.conf
/etc/ath.orig/apup
/etc/init.d/
/etc/init.d/madwifi-11n-scripts
/etc/rc.d/
/etc/rc.d/rc.wlan
/sbin/
/sbin/wlan_old
/sbin/update-wifi
/sbin/dni_apup
