# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/base-files-wnr1000v2-2.6.list
/
/sys/
/bin/
/bin/login
/bin/uci
/bin/ipcalc.sh
/rom/
/rom/note
/mnt
/usr/
/usr/bin/
/usr/lib/
/usr/lib/common.awk
/dev/
/lib/
/lib/config/
/lib/config/template.awk
/lib/config/uci.sh
/lib/config/parse_spec.awk
/lib/config/validate_config.awk
/lib/config/validate_spec.awk
/lib/config/specs/
/lib/config/specs/network.spec
/lib/config/validate.sh
/lib/config/uci.awk
/lib/firmware/
/lib/network/
/lib/network/config.sh
/proc/
/init
/tmp/
/www/
/etc/
/etc/preinit.arch
/etc/resolv.conf
/etc/profile
/etc/config/
/etc/config/system
/etc/config/network
/etc/rc.common
/etc/banner
/etc/hotplug.d/
/etc/hotplug.d/net/
/etc/hotplug.d/net/10-net
/etc/hotplug.d/iface/
/etc/hotplug.d/iface/00-netstate
/etc/hotplug.d/iface/10-routes
/etc/hotplug.d/usb/
/etc/hotplug.d/usb/10-usb-storage
/etc/functions.sh
/etc/inittab
/etc/preinit
/etc/crontabs/
/etc/mtab
/etc/passwd
/etc/shells
/etc/hotplug2-init.rules
/etc/group
/etc/ipkg.conf
/etc/protocols
/etc/init.d/
/etc/init.d/boot
/etc/init.d/detbutton
/etc/init.d/bridge
/etc/init.d/done
/etc/init.d/sysctl
/etc/init.d/rcS
/etc/sysctl.conf
/etc/hosts
/jffs/
/sbin/
/sbin/ifdown
/sbin/hotplug-call
/sbin/cmd_sipalg
/sbin/cmd_autoreset
/sbin/get-default-wan-mac-from-flash
/sbin/ifup
/sbin/run-ramfs
/sbin/qos.sh
/sbin/printhosts
/sbin/button_detecte
/sbin/usb-storage
/sbin/wifi
/var
