# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/busybox.list
/
/bin/
/bin/touch
/bin/adduser
/bin/uname
/bin/ln
/bin/chmod
/bin/ash
/bin/delgroup
/bin/iplink
/bin/addgroup
/bin/pwd
/bin/mv
/bin/ip
/bin/cat
/bin/nice
/bin/usleep
/bin/pidof
/bin/tar
/bin/cp
/bin/dd
/bin/sed
/bin/echo
/bin/busybox
/bin/iptunnel
/bin/sh
/bin/ping6
/bin/ping
/bin/grep
/bin/date
/bin/sync
/bin/chown
/bin/ls
/bin/ps
/bin/iproute
/bin/umount
/bin/egrep
/bin/rmdir
/bin/mknod
/bin/mount
/bin/su
/bin/sleep
/bin/fgrep
/bin/deluser
/bin/kill
/bin/rm
/bin/mkdir
/bin/df
/bin/ipaddr
/usr/
/usr/bin/
/usr/bin/[[
/usr/bin/awk
/usr/bin/killall5
/usr/bin/logger
/usr/bin/strings
/usr/bin/head
/usr/bin/killall
/usr/bin/cut
/usr/bin/expr
/usr/bin/md5sum
/usr/bin/wget
/usr/bin/passwd
/usr/bin/env
/usr/bin/[
/usr/bin/wc
/usr/bin/dirname
/usr/bin/printf
/usr/bin/hexdump
/usr/bin/tftp
/usr/bin/find
/usr/bin/test
/usr/bin/crontab
/usr/bin/free
/usr/sbin/
/usr/sbin/cmd_cron
/usr/sbin/crond
/usr/sbin/dsyslog
/etc/
/etc/bandcheck/
/etc/bandcheck/band-check
/etc/init.d/
/etc/init.d/cron
/sbin/
/sbin/halt
/sbin/fdisk
/sbin/klogd
/sbin/ifconfig
/sbin/switch_root
/sbin/sysctl
/sbin/route
/sbin/syslogd
/sbin/init
/sbin/insmod
/sbin/poweroff
/sbin/pivot_root
/sbin/logread
/sbin/reboot
/sbin/rmmod
/sbin/lsmod
