# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/dnsmasq.list
/
/usr/
/usr/sbin/
/usr/sbin/dnsmasq
/usr/sbin/dns-hijack
/etc/
/etc/dnsmasq.conf
/etc/init.d/
/etc/init.d/dnsmasq
