# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/igmpproxy.list
/
/etc/
/etc/init.d/
/etc/init.d/igmpproxy.init
/sbin/
/sbin/igmpproxy
/sbin/cmdigmp
