# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/kmod-ag7240-enet.list
/
/lib/
/lib/modules/
/lib/modules/2.6.15/
/lib/modules/2.6.15/ag7240_mod.ko
/etc/
/etc/modules.d/
/etc/modules.d/50-ag7240-enet
/etc/init.d/
/etc/init.d/net-lan
/etc/init.d/init-switch
/etc/init.d/net-wan
/sbin/
/sbin/ethreg
