# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/lltd.list
/
/usr/
/usr/sbin/
/usr/sbin/lld2d
/etc/
/etc/lld2d.conf
/etc/large.ico
/etc/init.d/
/etc/init.d/lltd
/etc/icon.ico
