# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/net-web.list
/
/usr/
/usr/sbin/
/usr/sbin/select_partition
/usr/sbin/wget_netgear
/usr/sbin/uhttpd
/usr/sbin/inetd
/www/
/www/AUTO_download_failure.htm
/www/UPG_upgrade.htm
/www/rbullet.gif
/www/port_fail.htm
/www/ipv6_pppoe.js
/www/lang_check.html
/www/router-info.htm
/www/BAS_basic.htm
/www/BRS_04_applySettings_wget.html
/www/BRS_03A_B_pppoe_reenter.html
/www/BRS_03A_A_noWan_check_net.html
/www/basic_home.htm
/www/top_nochange.html
/www/BKS_service_edit.htm
/www/AUTO_download_language.htm
/www/BRS_03A_E_IP_problem_staticIP.html
/www/UPG_cancel.htm
/www/AUTO_upgrade_no_langversion.htm
/www/RU_login_yes.htm
/www/reservation_add_wait.htm
/www/BAK_backup.htm
/www/bas_pppoe.js
/www/traffic.js
/www/RU_manual_fail.htm
/www/thank_login.htm
/www/CheckNewFirmware.htm
/www/welcomeok.htm
/www/FW_schedule.htm
/www/RU_login_no.htm
/www/UPG_version.htm
/www/BRS_04_applySettings.html
/www/BRS_01_checkNet.html
/www/pls_wait.html
/www/BAS_mulpppoe_add.htm
/www/RU_manu_local_resources.htm
/www/BAS_update.htm
/www/ping6_traceroute6_hidden_info.htm
/www/RU_isp_pppoe_static.htm
/www/ipv6.js
/www/BKS_keyword.htm
/www/script/
/www/script/brs_hijack_02_genieHelp.js
/www/script/brs_hijack_03A_A_noWan.js
/www/script/brs_genie_03A_B_pppoe_reenter.js
/www/script/brs_genie_03A_C_pptp_reenter.js
/www/script/brs_hijack_01_checkNet.js
/www/script/brs_hijack_03A_E_IP_problem.js
/www/script/ru_detcInetType.js
/www/script/brs_hijack_03B_haveBackupFile_fileRestore.js
/www/script/brs_check_manulConfig.js
/www/script/brs_hijack_04_applySettings_wget.js
/www/script/funcs.js
/www/script/brs_hijack_03A_C_pptp.js
/www/script/brs_hijack_03A_F_l2tp.js
/www/script/brs_hijack_03A_detcInetType.js
/www/script/brs_hijack_03A_D_bigpond.js
/www/script/brs_hijack_03A_E_IP_problem_staticIP_B_macClone.js
/www/script/brs_hijack_04_applySettings_ping.js
/www/script/brs_genie_03A_F_l2tp_reenter.js
/www/script/brs_hijack_03B_haveBackupFile_waitReboot.js
/www/script/brs_hijack_03A_B_pppoe.js
/www/script/brs_hijack_03B_haveBackupFile.js
/www/script/brs_hijack_05_networkIssue.js
/www/script/brs_hijack_03A_E_IP_problem_staticIP.js
/www/script/brs_hijack_index.js
/www/script/brs_hijack_success.js
/www/script/brs_04_B_checkNet.js
/www/script/brs_hijack_04_applySettings.js
/www/script/brs_hijack_03A_E_IP_problem_staticIP_A_inputIP.js
/www/script/brs_hijack_03A_E_IP_problem_staticIP_B_macClone_plzWait.js
/www/script/brs_hijack.js
/www/script/brs_hijack_wanlan_conflict.js
/www/script/brs_hiajck_03_haveBackupFile_ping.js
/www/script/brs_hijack_plzWait.js
/www/script/brs_hijack_03A_A_check_net.js
/www/RU_mac_spoof.js
/www/reservation.js
/www/WPS_client_status.htm
/www/BSW_wan.htm
/www/RU_isp_static.htm
/www/help.css
/www/ipv6_fixed.js
/www/CheckNewFW.html
/www/BRS_netgear_success.html
/www/BAK_factory.htm
/www/restore_process.htm
/www/IPv6_pppoe.htm
/www/FW_email.htm
/www/hidden_POT.htm
/www/RST_statistic.htm
/www/selects_disabled.gif
/www/upload.gif
/www/RST_conn_status.htm
/www/DEV_show_device_access.htm
/www/BRS_plzWait.html
/www/BRS_03A_C_pptp_reenter.html
/www/WLG_acl_add.htm
/www/UPG_process.htm
/www/hidden_info.htm
/www/AUTO_language_nowan.htm
/www/BRS_03A_E_IP_problem_staticIP_B_macClone.html
/www/AccessControl.js
/www/funcs.js
/www/UPNP_media.htm
/www/BRS_ping_result.html
/www/CheckNewFW_top.htm
/www/style.css
/www/enable_string_table.html
/www/BRS_03A_E_IP_problem.html
/www/BRS_04_applySettings_ping.html
/www/RU_detcInetType.html
/www/backup.js
/www/BRS_04_B_checkNet_ping.html
/www/BAS_mulpppoe.htm
/www/wlacl.js
/www/BRS_index.htm
/www/FW_log.htm
/www/forwarding.js
/www/WLG_adv.htm
/www/BKS_service_add.htm
/www/AUTO_write_language.htm
/www/help.htm
/www/BSW_wan_wait.htm
/www/debuginfo.htm
/www/UPnP_Media.js
/www/str_router.js
/www/BRS_retry.htm
/www/index.htm
/www/BAS_ether.htm
/www/BRS_03B_haveBackupFile_ping.html
/www/ftp_no_shares.htm
/www/LAN_lan.htm
/www/RU_isp_l2tp.htm
/www/adv_index.htm
/www/RU_isp_list.js
/www/wlan.js
/www/WPS_status.htm
/www/upg_get_status.htm
/www/upgrade.js
/www/h_BAS_mulpppoe.html
/www/UPG_failure.htm
/www/BRS_03A_F_l2tp_reenter.html
/www/BRS_03B_haveBackupFile.html
/www/BRS_wanlan_conflict.html
/www/FW_forward.htm
/www/WPS_pin.htm
/www/reservation_edit.htm
/www/UPNP_upnp.htm
/www/RU_finish_manual.htm
/www/WLG_acl.htm
/www/wan.js
/www/ftp_via_no_shares.htm
/www/UPG_cancel_hidden.htm
/www/LGO_logout.htm
/www/BRS_success.html
/www/triggering.js
/www/remote_fail.htm
/www/CheckNewFW_upgrade.htm
/www/BAS_update_wan.htm
/www/reboot.htm
/www/multi_login.html
/www/BRS_check_manulConfig.html
/www/IPv6_auto.htm
/www/AUTO_upgrade_no_internet.htm
/www/forwarding.htm
/www/RU_manual_pppoe_dual.htm
/www/change_to_wireless.htm
/www/liteblue.gif
/www/AUTO_start.htm
/www/AUTO_upgrade_noversion.htm
/www/schedule.js
/www/reboot_waiting.htm
/www/BSW_wan_query.htm
/www/del_fail.htm
/www/forwarding_edit.htm
/www/AccessControl_show.htm
/www/upnp.js
/www/StringTableUpload.html
/www/center_language.html
/www/top.js
/www/BRS_03A_E_IP_problem_staticIP_A_inputIP.html
/www/checkUpdate_language.htm
/www/image/
/www/image/close_list.gif
/www/image/basic_homes.gif
/www/image/buttons.gif
/www/image/ca_purple_long_button.gif
/www/image/advanced_homes.gif
/www/image/process_bar.jpg
/www/image/top-oneline.gif
/www/image/ca-little-green-buitton.gif
/www/image/wl_ap.gif
/www/image/wait30.gif
/www/image/tops.gif
/www/image/ca-little-purple-button.gif
/www/image/top_button_background.gif
/www/image/help/
/www/image/help/help_up.gif
/www/image/help/help_down.gif
/www/image/help/help-bar-icon.gif
/www/image/help/help-bar.gif
/www/image/menus.gif
/www/image/footer/
/www/image/footer/search_button.gif
/www/image/footer/search_button_middle.gif
/www/image/footer/footer_doub.gif
/www/image/footer/search_button_long.gif
/www/image/footer/footer.gif
/www/image/ca_text_pic.gif
/www/image/open_list.gif
/www/image/spacer.gif
/www/WPS.htm
/www/RU_manual_pppoe.htm
/www/top.html
/www/RU_success.htm
/www/BRS_02_genieHelp.html
/www/STR_routes.htm
/www/pls_wait_factory_reboot.html
/www/RST_status.htm
/www/BRS_03A_C_pptp.html
/www/basic.js
/www/BAS_buttom.htm
/www/Add_WPS_Client.htm
/www/basic_wait.htm
/www/wlan_adv.js
/www/pls_wait_updateLanguage.html
/www/forwarding_add.htm
/www/RU_isp_pptp.htm
/www/POT.htm
/www/STR_routes_add.htm
/www/NoInternetConn.html
/www/DEV_device.htm
/www/RU_manual_l2tp.htm
/www/region_class.js
/www/show_traffic.htm
/www/start.htm
/www/goodbye.html
/www/bas_ether.js
/www/FW_remote.htm
/www/WIZ_sel.htm
/www/BRS_top.html
/www/wps_client_status.html
/www/IPv6_autoConfig.htm
/www/RU_welcome.htm
/www/securityquestions.cgi
/www/BRS_04_B_checkNet.html
/www/RU_manual_pptp.htm
/www/WLG_schedule_an_edit.htm
/www/RU_isp_spoof.htm
/www/RU_pppoe_dual.htm
/www/WLG_acl_wait.htm
/www/selects.gif
/www/hidden_channel_wifi_test.htm
/www/WPS_failure.htm
/www/BRS_03A_E_IP_problem_staticIP_B_macClone_plzWait.html
/www/logs.js
/www/BRS_03A_B_pppoe.html
/www/lang_search.htm
/www/wps_icon_on.jpg
/www/AUTO_search.htm
/www/forwarding_add_wait.htm
/www/BKS_service.htm
/www/BAS_mulpppoe_edit.htm
/www/language/
/www/language/Chinese.js
/www/language/Italian.js
/www/language/Swedish.js
/www/language/German.js
/www/language/Portuguese.js
/www/language/Japanese.js
/www/language/Romanian.js
/www/language/Spanish.js
/www/language/Slovak.js
/www/language/Finnish.js
/www/language/Russian.js
/www/language/Norwegian.js
/www/language/Polish.js
/www/language/Korean.js
/www/language/Dutch.js
/www/language/Danish.js
/www/language/French.js
/www/language/RU_flow.js
/www/language/Turkish.js
/www/language/Czech.js
/www/language/English.js
/www/language/Slovenian.js
/www/language/RU_flow_en.js
/www/language/Hungarian.js
/www/language/Greek.js
/www/PWD_password.htm
/www/continue.htm
/www/BRS_03B_haveBackupFile_fileRestore.html
/www/passwd.js
/www/BRS_03A_F_l2tp.html
/www/WPS_success.htm
/www/WLG_schedule_an_add.htm
/www/RST_status_wds.htm
/www/cgi-bin/
/www/cgi-bin/firewall.sh
/www/cgi-bin/uhttpd.init
/www/cgi-bin/firewall_function.sh
/www/cgi-bin/uhttpd.sh
/www/RU_manual_spoof.htm
/www/top_submit.html
/www/IPv6_index.htm
/www/forwarding_edit_wait.htm
/www/bas_mulpppoe.js
/www/hidden_manual_set_ant.html
/www/DNS_ddns.htm
/www/pls_wait_qos_check.html
/www/IPv6_passThrougth.htm
/www/pls_wait_reboot.html
/www/IPv6_tunnel.htm
/www/lan.js
/www/help/
/www/help/_qos_help.html
/www/help/_FW_forward_help.html
/www/help/_Wizard_setup_help.html
/www/help/_BKS_service_edit_help.html
/www/help/_IPv6_fixed_help.html
/www/help/_FW_trigger_help.html
/www/help/_wps_help.html
/www/help/_upgrade_help.html
/www/help/_upnp_help.html
/www/help/_backup_help.html
/www/help/_wds_help.html
/www/help/_IPv6_autoConfig_help.html
/www/help/_BKS_servic_help.html
/www/help/_USB_advance_help.html
/www/help/_BKS_service_help.html
/www/help/_access_control_help.html
/www/help/_plc_avail_dev_help.html
/www/help/_IPv6_dhcp_help.html
/www/help/_apmode_help.html
/www/help/_BKS_keyword_help.html
/www/help/_remote_share_help.html
/www/help/help.htm
/www/help/_wlan_guest_bgn_help.html
/www/help/_BAS_ether_help.html
/www/help/_WLG_acl_help.html
/www/help/_FW_email_help.html
/www/help/_BAS_pppoe_help.html
/www/help/_UPNP_media_help.html
/www/help/_IPv6_tunnel_help.html
/www/help/_wlan_guest_an_help.html
/www/help/_IPv6_help.html
/www/help/_BKS_service_add_help.html
/www/help/_USB_basic_help.html
/www/help/_info_help.html
/www/help/_FW_remote_help.html
/www/help/_lan_help.html
/www/help/_WLG_acl_add_help.html
/www/help/_wladv_help.html
/www/help/_IPv6_auto_help.html
/www/help/_FW_schedule_help.html
/www/help/_attach_device_help.html
/www/help/_DNS_ddns_help.html
/www/help/_passwd_help.html
/www/help/_USB_approve_help.html
/www/help/_Wizard_cfm_help.html
/www/help/_CF_log_help.html
/www/help/_wlan_help.html
/www/help/_STR_routes_help.html
/www/help/_BAS_pptp_help.html
/www/help/_traffic_help.html
/www/help/_IPv6_pppoe_help.html
/www/help/_print_share_help.html
/www/help/_WAN_wan_help.html
/www/help/_USB_setting_help.html
/www/BRS_03A_A_noWan.html
/www/WAN_wan.htm
/www/currentsetting.htm
/www/DNS_status_wait.htm
/www/RU_welcome.js
/www/AUTO_download_image.htm
/www/UPG_check_version.htm
/www/mp_lang.js
/www/BRS_03A_detcInetType.html
/www/RU_isp_l2tp_static.htm
/www/RU_MC_login.htm
/www/advanced.js
/www/DEV_show_device.htm
/www/ipv6_autoConfig.js
/www/BRS_04_applySettings_wgetResult.html
/www/DNS_oray_status.htm
/www/block_sites.js
/www/darkblue.gif
/www/UPG_check_version_hidden.htm
/www/IPv6_dhcp.htm
/www/IPv6_disabled.htm
/www/IPv6_fixed.htm
/www/DNS_ddns_status.htm
/www/check_passwd.htm
/www/bas_pptp.js
/www/access_control.htm
/www/wadv_sechdule.js
/www/triggering.htm
/www/RU_check_wan_port.html
/www/BRS_01_checkNet_ping.html
/www/RU_WIZ_sel.htm
/www/BAS_pppoe.htm
/www/change_domain.htm
/www/wps_icon_off.jpg
/www/block_services.js
/www/BRS_03B_haveBackupFile_change_domain.htm
/www/triggering_add.htm
/www/email.js
/www/ipv6_tunnel.js
/www/RU_isp_pptp_static.htm
/www/WLG_acl_edit.htm
/www/WLG_schedule_edit.htm
/www/WPS_client.htm
/www/reservation_add.htm
/www/AUTO_upgrade_process.htm
/www/ap_waiting.htm
/www/STR_routes_edit.htm
/www/BRS_05_networkIssue.html
/www/style/
/www/style/basic.css
/www/style/hijack_05_networkIssue.css
/www/style/hijack_03A_E_IP_problem_staticIP.css
/www/style/help.css
/www/style/hijack_03B_haveBackupFile.css
/www/style/hijack_03A_E_IP_problem_staticIP_A_inputIP.css
/www/style/hijack_03A_C_pptp.css
/www/style/internet.css
/www/style/form.css
/www/style/top.css
/www/style/hijack_03A_D_bigpond.css
/www/style/hijack_ru_finish_manual.css
/www/style/basic_home.css
/www/style/hijack_ru_welcome.css
/www/style/hijack_03A_E_IP_problem.css
/www/style/advanced_home.css
/www/style/top_style.css
/www/style/hijack_03A_C_pptp_reenter.css
/www/style/hijack_success.css
/www/style/hijack_03A_wanInput_reenter.css
/www/style/hijack_03A_E_IP_problem_staticIP_B_macClone.css
/www/style/hijack_style.css
/www/style/advanced.css
/www/style/hijack_02_genieHelp.css
/www/style/hijack_03A_wanInput.css
/www/style/hijack_ru_isp_fail.css
/www/style/advanced_home_ru.css
/www/WLG_schedule_close.htm
/www/triggering_edit.htm
/www/ipv6_dhcp.js
/www/AUTO_MD5_failure.htm
/www/ddns.js
/www/RU_isp_pppoe.htm
/www/remote.js
/www/WLG_wireless.htm
/www/WLG_schedule_add.htm
/www/edit_fail.htm
/www/BRS_03B_haveBackupFile_waitReboot.html
/www/languages-en.js
/www/traffic.htm
/www/RU_manual_l2tp_dual.htm
/www/wps_get_status.htm
/www/RU_manual_pptp_dual.htm
/www/RU_isp_fail.htm
/www/process_bar.gif
/www/BAS_pptp.htm
/www/RU_manual_pppoe2.htm
/www/spacer.gif
/www/hidden_channel_wifi_coexist.htm
/www/BRS_03A_D_bigpond.html
/etc/
/etc/init.d/
/etc/init.d/uhttpd.init
