# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/ntpclient.list
/
/usr/
/usr/sbin/
/usr/sbin/ntpclient
/etc/
/etc/init.d/
/etc/init.d/ntpclient
