# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/trafficmeter.list
/
/usr/
/usr/traffic_meter/
/usr/traffic_meter/traffic_meter
/usr/traffic_meter/traffic_statistics
/usr/traffic_meter/ppp_last_conn_time
/usr/traffic_meter/network_status
/usr/traffic_meter/ppp_last_stop_time
/etc/
/etc/init.d/
/etc/init.d/traffic_meter.init
/sbin/
/sbin/traffic_meter
/sbin/cmd_traffic_meter
/sbin/tfm_led
/sbin/generate_traffic_meter_conf
