# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/utelnetd.list
/
/usr/
/usr/sbin/
/usr/sbin/utelnetd
/etc/
/etc/init.d/
/etc/init.d/utelnetd
