# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/wlan-common.list
/
/etc/
/etc/config/
/etc/config/wireless
/etc/hotplug.d/
/etc/hotplug.d/wps/
/etc/hotplug.d/wps/00-wps
/etc/hotplug.d/wps/01-wps-led
/etc/init.d/
/etc/init.d/wlan-common
/sbin/
/sbin/wpsled
/sbin/wlan
