# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/wlanlog.list
/
/usr/
/usr/bin/
/usr/bin/wlanlog
/etc/
/etc/init.d/
/etc/init.d/wlanlog
