# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/usr/lib/ipkg/info/zebra.list
/
/usr/
/usr/sbin/
/usr/sbin/ripngd
/usr/sbin/zebra
/etc/
/etc/ripngd.conf
/etc/init.d/
/etc/init.d/ripngd
/etc/init.d/zebra
/etc/zebra.conf
