# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K stop
tts/0::respawn:/bin/ash --login
ttyS0::respawn:/bin/ash --login
tty1::respawn:/bin/ash --login
