# firmware/extracted/_WNR1000v2_V1_1_2_60NA_img_extracted/squashfs-root/etc/hotplug.d/usb/10-usb-storage
#!/bin/sh

# Copyright (C) 2006 OpenWrt.org
echo "called with DEVPATH=$DEVPATH, INTERFACE=$INTERFACE"
case "$ACTION" in
	add)
		[ -n "${INTERFACE}" ] &&
		[ "$(expr substr ${INTERFACE} 1 2)" == "8/" ]  && {
			/sbin/usb-storage &
		}
		;;
	remove)
		ls /sys$DEVPATH
		[ -n "${INTERFACE}" ] &&
		[ "$(expr substr ${INTERFACE} 1 2)" == "8/" ]  && {
			/sbin/usb-storage.remove &
		}
		;;
esac	
		
